# sesrv-timer-1.service
[Unit]
Description=SeSrv Script Timer 1 Service

[Service]
Type=oneshot
ExecStart=/usr/bin/sesrv-script timer_one
